# service/cpupowerperf.timer
[Unit]
Description=CPU Power Performance

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
